# supervisor/cron.conf
[program:cron]
command=cron -f -L 15
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/cron.log
